# deployments/systemd/multinic-agent.service
[Unit]
Description=MultiNIC Network Interface Agent
Documentation=https://github.com/Minzino/multinic-agent
After=network-online.target mysql.service mariadb.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/multinic-agent
Restart=always
RestartSec=30s

# 환경 변수 설정
Environment="LOG_LEVEL=info"
Environment="POLL_INTERVAL=30s"
Environment="HEALTH_PORT=8080"
EnvironmentFile=-/etc/multinic-agent/config.env

# 보안 설정 (systemd가 관리)
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=multinic-agent

# 리소스 제한
LimitNOFILE=65536
LimitNPROC=4096

# 재시작 정책
StartLimitInterval=600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target